# testing/docker/builder-bmoweb/supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/run-apache.sh
autorestart = true
redirect_stderr = true
