# autoland/upstart/autoland.conf
description "Autoland"
author "auto-tools@mozilla.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

setuid ubuntu
setgid ubuntu

script
  cd /home/ubuntu/version-control-tools/autoland/autoland
  exec python autoland.py --log-path=/home/ubuntu/autoland.log
end script
