# ansible/roles/hg-ssh/files/snsnotifier.service
[Unit]
Description=Notify SNS about Mercurial events
After=network.target remote-fs.target nss-lookup.target
PartOf=hg-master.target
OnFailure=mail-unit-output@%n.service

# We only want to have a single master at a time. Look for a file with the
# current hostname in it indicating that we are the master.
AssertPathExists=/repo/hg/master.%H

[Service]
WorkingDirectory=/
User=hg-notifier
Group=hg-notifier
PrivateTmp=true

ExecStart=/var/hg/venv_tools/bin/vcsreplicator-sns-notifier /etc/mercurial/notifications.ini
Restart=always
TimeoutStopSec=60

# Needed for HTTPS_PROXY in some environments.
EnvironmentFile=/etc/environment

SyslogIdentifier=snsnotifier
SyslogFacility=local2

[Install]
WantedBy=hg-master.target
