# ansible/roles/hg-web/files/vcsreplicator@.service
[Unit]
Description=Mirror Mercurial changes
After=network.target remote-fs.target nss-lookup.target

[Service]
ExecStart=/var/hg/venv_replication/bin/vcsreplicator-consumer /etc/mercurial/vcsreplicator.ini --partition %i
WorkingDirectory=/
PrivateTmp=true
User=hg
Group=hg

# We want new repos to be group writable by default. Default umask of
# 022 prevents this.
UMask=002

# Transient errors are absorbed by `retry_on_failure`; the burst limit
# below trips `failed` state on rapid crash loops. `StartLimitInterval`
# (legacy alias) is required for systemd 219 / RHEL 7.
StartLimitInterval=300
StartLimitBurst=20
Restart=always
RestartSec=3
TimeoutStopSec=60

SyslogIdentifier=vcsreplicator
SyslogFacility=local2

[Install]
WantedBy=multi-user.target
