# ansible/roles/kafka-broker/files/kafka-logs-cleanup.timer
[Unit]
Description=Cleanup kafka/zookeeper log files

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
